# agent-viz-v2.service
[Unit]
Description=AgentViz v2 - Real-time Agent Activity Visualization Dashboard
Documentation=https://github.com/setrox/agent-viz-v2
After=network.target

[Service]
Type=simple
User=setrox
WorkingDirectory=/home/setrox/.openclaw/agent-viz-v2
ExecStart=/usr/bin/npm run server
Environment=NODE_ENV=production
Environment=WS_PORT=18789
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/setrox/.openclaw/agent-viz-v2

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=default.target
